# scripts/com.agentarchive.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.agentarchive.sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mac/agent-archive/.venv/bin/python</string>
    <string>-m</string>
    <string>agent_archive.cli</string>
    <string>sync</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict><key>AGENT_ARCHIVE_ROOT</key><string>/Users/mac/agent-archive-data</string></dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>/Users/mac/agent-archive-data/sync.log</string>
  <key>StandardErrorPath</key><string>/Users/mac/agent-archive-data/sync.err</string>
  <key>WorkingDirectory</key><string>/Users/mac/agent-archive</string>
</dict>
</plist>
