# scripts/com.agentarchive.digest-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.agentarchive.digest-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mac/agent-archive/scripts/push-digest.sh</string>
    <string>day</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_ARCHIVE_ROOT</key><string>/Users/mac/agent-archive-data</string>
    <key>AGENT_ARCHIVE_LARK_CHAT_ID</key><string>__CHAT_ID__</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>40</integer></dict>
  <key>StandardOutPath</key><string>/Users/mac/agent-archive-data/digest.log</string>
  <key>StandardErrorPath</key><string>/Users/mac/agent-archive-data/digest.err</string>
</dict>
</plist>
